# frameworks-webhook.service
[Unit]
Description=Frameworks to Klaviyo Webhook Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/frameworks-to-klaviyo
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/home/ubuntu/frameworks-to-klaviyo/.env

[Install]
WantedBy=multi-user.target
